# ansible/roles/jenkins/templates/jenkins.service.j2
# ============================================
# Jenkins systemd service override (optional)
# Currently not used, kept for future tuning
# ============================================

[Service]
Environment="JAVA_OPTS=-Djava.awt.headless=true"
